# deploy/supervisor/aulog.conf
; AuLog Web 服务 — Supervisor 配置
;
; 安装步骤（在 ECS 上执行）：
;   sudo apt install -y supervisor
;   sudo cp /root/AuLog/deploy/supervisor/aulog.conf /etc/supervisor/conf.d/aulog.conf
;   sudo mkdir -p /root/AuLog/logs
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start aulog
;
; 常用命令：
;   sudo supervisorctl status aulog
;   sudo supervisorctl restart aulog
;   sudo supervisorctl stop aulog
;   tail -f /root/AuLog/logs/aulog.stdout.log

[program:aulog]
command=/root/AuLog/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 18888
directory=/root/AuLog
user=root
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
priority=999

; 应用通过 app/db.py 的 load_dotenv() 读取 /root/AuLog/.env
environment=PATH="/root/AuLog/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

stdout_logfile=/root/AuLog/logs/aulog.stdout.log
stderr_logfile=/root/AuLog/logs/aulog.stderr.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
redirect_stderr=false
